# ai-services/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:multimedia-api]
command=python /app/ai-services/multimedia_api_v4.py
directory=/app/ai-services
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/multimedia-api.err.log
stdout_logfile=/var/log/supervisor/multimedia-api.out.log
environment=PYTHONPATH="/app/ai-services"

[program:enhanced-indexer]
command=python /app/ai-services/enhanced_indexer_v4.py
directory=/app/ai-services
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/indexer.err.log
stdout_logfile=/var/log/supervisor/indexer.out.log
environment=PYTHONPATH="/app/ai-services"

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface